# clickhouse_hass/rootfs/etc/s6-overlay/s6-rc.d/init-clickhouse-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init clickhouse-server - generate config from options
# ==============================================================================

set -e

bashio::log.info "Running init script..."

if [ ! -d "/data/config/clickhouse-server" ]; then
    bashio::log.info "Create dir for configurations clickhouse-server..."
    mkdir -p /data/config/
    cp -r /opt/etc/clickhouse-server /data/config
fi

if [ ! -d "/data/config/clickhouse-client" ]; then
    bashio::log.info "Create dir for configurations clickhouse-client..."
    mkdir -p /data/config/clickhouse-client
    cp -r /opt/etc/clickhouse-client /data/config/
fi

if [ ! -d "/data/clickhouse/data" ]; then
    bashio::log.info "Create dir for data clickhouse-server..."
    mkdir -p /data/clickhouse/data
fi

# ── Ports ────────────────────────────────────────────────────────────────────
HTTP_PORT=""
TCP_PORT=""
MYSQL_PORT=""
POSTGRESQL_PORT=""

if bashio::config.has_value 'http_port'; then
    HTTP_PORT="<http_port>$(bashio::config 'http_port')</http_port>"
fi
if bashio::config.has_value 'tcp_port'; then
    TCP_PORT="<tcp_port>$(bashio::config 'tcp_port')</tcp_port>"
fi
if bashio::config.has_value 'mysql_port'; then
    MYSQL_PORT="<mysql_port>$(bashio::config 'mysql_port')</mysql_port>"
fi
if bashio::config.has_value 'postgresql_port'; then
    POSTGRESQL_PORT="<postgresql_port>$(bashio::config 'postgresql_port')</postgresql_port>"
fi

# ── Logging ──────────────────────────────────────────────────────────────────
LOG_LEVEL="information"
if bashio::config.has_value 'logging.level'; then
    LOG_LEVEL="$(bashio::config 'logging.level')"
fi

LOG_QUERIES=""
if bashio::config.has_value 'logging.log_queries'; then
    if bashio::var.true "$(bashio::config 'logging.log_queries')"; then
        LOG_QUERIES="<log_queries>1</log_queries>"
    else
        LOG_QUERIES="<log_queries>0</log_queries>"
    fi
fi

# ── Performance ───────────────────────────────────────────────────────────────
MAX_CONNECTIONS=""
MAX_CONCURRENT_QUERIES=""
MAX_SERVER_MEMORY_USAGE=""
MAX_SERVER_MEMORY_USAGE_TO_RAM_RATIO=""
BACKGROUND_POOL_SIZE=""

if bashio::config.has_value 'performance.max_connections'; then
    MAX_CONNECTIONS="<max_connections>$(bashio::config 'performance.max_connections')</max_connections>"
fi
if bashio::config.has_value 'performance.max_concurrent_queries'; then
    MAX_CONCURRENT_QUERIES="<max_concurrent_queries>$(bashio::config 'performance.max_concurrent_queries')</max_concurrent_queries>"
fi
if bashio::config.has_value 'performance.max_server_memory_usage'; then
    MAX_SERVER_MEMORY_USAGE="<max_server_memory_usage>$(bashio::config 'performance.max_server_memory_usage')</max_server_memory_usage>"
fi
if bashio::config.has_value 'performance.max_server_memory_usage_to_ram_ratio'; then
    MAX_SERVER_MEMORY_USAGE_TO_RAM_RATIO="<max_server_memory_usage_to_ram_ratio>$(bashio::config 'performance.max_server_memory_usage_to_ram_ratio')</max_server_memory_usage_to_ram_ratio>"
fi
if bashio::config.has_value 'performance.background_pool_size'; then
    BACKGROUND_POOL_SIZE="<background_pool_size>$(bashio::config 'performance.background_pool_size')</background_pool_size>"
fi

# ── Security ─────────────────────────────────────────────────────────────────
DEFAULT_PASSWORD=""
if bashio::config.has_value 'security.default_password'; then
    DEFAULT_PASSWORD="$(bashio::config 'security.default_password')"
fi

# ── Storage ───────────────────────────────────────────────────────────────────
MAX_TABLE_SIZE_TO_DROP=""
MAX_PARTITION_SIZE_TO_DROP=""

if bashio::config.has_value 'storage.max_table_size_to_drop'; then
    MAX_TABLE_SIZE_TO_DROP="<max_table_size_to_drop>$(bashio::config 'storage.max_table_size_to_drop')</max_table_size_to_drop>"
fi
if bashio::config.has_value 'storage.max_partition_size_to_drop'; then
    MAX_PARTITION_SIZE_TO_DROP="<max_partition_size_to_drop>$(bashio::config 'storage.max_partition_size_to_drop')</max_partition_size_to_drop>"
fi

# ── Generate hass.xml ─────────────────────────────────────────────────────────
mkdir -p /data/config/clickhouse-server/config.d
cat > /data/config/clickhouse-server/config.d/hass.xml << EOF
<?xml version="1.0" encoding="UTF-8"?>
<clickhouse>
    <logger>
        <level>${LOG_LEVEL}</level>
        <console>1</console>
        ${LOG_QUERIES}
    </logger>
    <path>/data/clickhouse/data</path>
    <tmp_path>/data/clickhouse/data/tmp/</tmp_path>
    <custom_cached_disks_base_directory>/data/clickhouse/data/caches/</custom_cached_disks_base_directory>
    <user_files_path>/data/clickhouse/data/user_files/</user_files_path>
    <format_schema_path>/data/clickhouse/data/format_schemas/</format_schema_path>
    <user_directories>
        <local_directory>
            <path>/data/clickhouse/data/access/</path>
        </local_directory>
    </user_directories>
    <listen_host>0.0.0.0</listen_host>
    ${HTTP_PORT}
    ${TCP_PORT}
    ${MYSQL_PORT}
    ${POSTGRESQL_PORT}
    ${MAX_CONNECTIONS}
    ${MAX_CONCURRENT_QUERIES}
    ${MAX_SERVER_MEMORY_USAGE}
    ${MAX_SERVER_MEMORY_USAGE_TO_RAM_RATIO}
    ${BACKGROUND_POOL_SIZE}
    ${MAX_TABLE_SIZE_TO_DROP}
    ${MAX_PARTITION_SIZE_TO_DROP}
</clickhouse>
EOF

sed -i '/^[[:space:]]*$/d' /data/config/clickhouse-server/config.d/hass.xml

# ── Default user password ─────────────────────────────────────────────────────
if [ -n "${DEFAULT_PASSWORD}" ]; then
    PASSWORD_SHA256=$(echo -n "${DEFAULT_PASSWORD}" | sha256sum | awk '{print $1}')
    cat > /data/config/clickhouse-server/config.d/hass-users.xml << EOF
<?xml version="1.0" encoding="UTF-8"?>
<clickhouse>
    <users>
        <default>
            <password_sha256_hex>${PASSWORD_SHA256}</password_sha256_hex>
        </default>
    </users>
</clickhouse>
EOF
    bashio::log.info "Default user password set."
else
    rm -f /data/config/clickhouse-server/config.d/hass-users.xml
fi

chown clickhouse:clickhouse /data/clickhouse/data/
chmod ugo+Xrw -R /data/clickhouse/data/ /data/config/clickhouse-server/ /data/config/clickhouse-client/

bashio::log.info "Init script done!"
