# clickhouse_hass/rootfs/etc/s6-overlay/s6-rc.d/init-clickhouse-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init clickhouse-server - generate config from options
# ==============================================================================

set -e

bashio::log.info "Running init script..."

if [ ! -d "/data/config/clickhouse-server" ]; then
    bashio::log.info "Create dir for configurations clickhouse-server..."
    mkdir -p /data/config/
    cp -r /opt/etc/clickhouse-server /data/config
fi

if [ ! -d "/data/config/clickhouse-client" ]; then
    bashio::log.info "Create dir for configurations clickhouse-client..."
    mkdir -p /data/config/clickhouse-client
    cp -r /opt/etc/clickhouse-client /data/config/
fi

if [ ! -d "/data/clickhouse/data" ]; then
    bashio::log.info "Create dir for data clickhouse-server..."
    mkdir -p /data/clickhouse/data
fi

# Regenerate hass.xml from scratch on every startup to avoid duplicate elements
mkdir -p /data/config/clickhouse-server/config.d
cat > /data/config/clickhouse-server/config.d/hass.xml << 'EOF'
<?xml version="1.0" encoding="UTF-8"?>
<clickhouse>
    <logger>
        <level>information</level>
        <console>1</console>
    </logger>
    <path>/data/clickhouse/data</path>
    <tmp_path>/data/clickhouse/data/tmp/</tmp_path>
    <custom_cached_disks_base_directory>/data/clickhouse/data/caches/</custom_cached_disks_base_directory>
    <user_files_path>/data/clickhouse/data/user_files/</user_files_path>
    <format_schema_path>/data/clickhouse/data/format_schemas/</format_schema_path>
    <user_directories>
        <local_directory>
            <path>/data/clickhouse/data/access/</path>
        </local_directory>
    </user_directories>
    <listen_host>0.0.0.0</listen_host>
    <http_port>8124</http_port>
    <tcp_port>9000</tcp_port>
    <mysql_port>9004</mysql_port>
    <postgresql_port>9005</postgresql_port>
</clickhouse>
EOF

chown clickhouse:clickhouse /data/clickhouse/data/
chmod ugo+Xrw -R /data/clickhouse/data/ /data/config/clickhouse-server/ /data/config/clickhouse-client/

bashio::log.info "Init script done!"
