# clickhouse_hass/rootfs/etc/s6-overlay/s6-rc.d/clickhouse-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start clickhouse-server
# ==============================================================================

set -e

bashio::log.info "Starting clickhouse-server..."
exec s6-setuidgid clickhouse /usr/bin/clickhouse-server --config-file=/data/config/clickhouse-server/config.xml
