# rabbitmq_hass/rootfs/etc/s6-overlay/s6-rc.d/init-rabbitmq/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init rabbitmq - generate config from options
# ==============================================================================

set -e

bashio::log.level "$(bashio::config 'log_level')"
bashio::log.info "Running init script..."

RABBITMQ_CONFIG="/data/config"
RABBITMQ_CONFIG_EXTEND="${RABBITMQ_CONFIG}/conf.d"
RABBITMQ_DATA_DIR="/data/rabbitmq"
RABBITMQ_LOG_DIR="/var/log/rabbitmq"

# ── First-run: copy default configs if not present ───────────────────────────
if [ ! -f "${RABBITMQ_CONFIG}/rabbitmq.conf" ]; then
    bashio::log.info "Copying default rabbitmq.conf..."
    mkdir -p "${RABBITMQ_CONFIG}"
    cp /opt/rabbitmq/rabbitmq.conf "${RABBITMQ_CONFIG}/rabbitmq.conf"
fi

if [ ! -f "${RABBITMQ_CONFIG}/definitions.json" ]; then
    bashio::log.info "Copying default definitions.json..."
    mkdir -p "${RABBITMQ_CONFIG}"
    cp /opt/rabbitmq/definitions.json "${RABBITMQ_CONFIG}/definitions.json"
fi

if [ ! -d "${RABBITMQ_CONFIG_EXTEND}" ]; then
    bashio::log.info "Copying default conf.d/..."
    mkdir -p "${RABBITMQ_CONFIG_EXTEND}"
    cp /opt/rabbitmq/conf.d/* "${RABBITMQ_CONFIG_EXTEND}/"
fi

if [ ! -d "${RABBITMQ_DATA_DIR}" ]; then
    bashio::log.info "Creating data directory..."
    mkdir -p "${RABBITMQ_DATA_DIR}"
fi

# ── Apply log level from HA config ───────────────────────────────────────────
LOG_LEVEL=$(bashio::config 'log_level')
sed -i "s|^log.console.level = .*|log.console.level = ${LOG_LEVEL}|" "${RABBITMQ_CONFIG}/rabbitmq.conf"

# ── Enable plugins ────────────────────────────────────────────────────────────
plugins=$(bashio::config 'plugins' | jq -r 'to_entries | .[] | select(.value==true) | .key' | paste -sd "," -)
bashio::log.info "Enabling plugins: ${plugins:-none}"
echo "[${plugins}]." > "${RABBITMQ_CONFIG}/enabled_plugins"

# ── Set correct permissions ───────────────────────────────────────────────────
chown -R rabbitmq:rabbitmq "${RABBITMQ_CONFIG}" "${RABBITMQ_DATA_DIR}" "${RABBITMQ_LOG_DIR}"

bashio::log.info "Init script completed successfully!"
