# postfix_hass/rootfs/etc/s6-overlay/s6-rc.d/syslog-ng/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start syslog-ng service
# ==============================================================================

set -e

bashio::log.info "Running syslog-ng..."
/usr/sbin/syslog-ng --cfgfile=/etc/syslog-ng/syslog-ng.conf -F --persist-file /run/syslog-ng.persist --control=/run/syslog-ng.ctl

#exec 2>&1 s6-setuidgid abc /usr/sbin/syslog-ng -F -f /etc/syslog-ng/syslog-ng.conf --persist-file /run/syslog-ng.persist --pidfile=/run/syslog-ng.pid --control=/run/syslog-ng.ctl --stderr --no-caps
bashio::log.info "Running syslog-ng... Done"
