# postfix_hass/rootfs/etc/s6-overlay/s6-rc.d/init-postfix/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init postfix - generate config from options
# ==============================================================================

set -e

bashio::log.info "Running init script..."

bashio::log.info "Merging options & variables for template"
# shellcheck disable=SC2046
JSON_CONF=$(jq '({options: .})' /data/options.json)

bashio::log.info "Generating postfix/main.cf from template in /etc/postfix/main.cf.gtpl"
# shellcheck disable=SC2086
echo $JSON_CONF | tempio \
    -template /etc/postfix/main.cf.gtpl \
    -out /etc/postfix/main.cf.tmp

bashio::log.info "Init script done!"