# postfix_hass/rootfs/etc/s6-overlay/s6-rc.d/init-postfix/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init postfix - generate config from options
# ==============================================================================

set -e

bashio::log.info "Running init script..."

if ! DOMAIN=$(bashio::config 'myhostname'); then
    bashio::log.error "Missing 'myhostname' configuration."
    exit 1
fi

if ! mynetworks=$(bashio::config 'mynetworks' | jq -r 'join(", ")'); then
    bashio::log.error "Missing or invalid 'mynetworks' configuration."
    exit 1
fi

# add domain
bashio::log.info "Set hostname to $DOMAIN..."
postconf -e myhostname="$DOMAIN"
postconf -e mydestination="localhost"
echo "$DOMAIN" > /etc/mailname
echo "$DOMAIN" > /etc/hostname
bashio::log.info "Set hostname to $DOMAIN... Done"

bashio::log.info "Merging options & variables for template"
# shellcheck disable=SC2046
JSON_CONF=$(jq '({options: .})' /data/options.json)

bashio::log.info "Generating postfix/main.cf from template in /etc/postfix/main.cf.gtpl"
# shellcheck disable=SC2086
echo $JSON_CONF | tempio \
    -template /etc/postfix/main.cf.gtpl \
    -out /etc/postfix/main.cf.tmp

# add mynetworks in config
bashio::log.info "Configure mynetworks: ${mynetworks}"

if grep -q "^mynetworks = " /etc/postfix/main.cf; then
    sed -i "s|^mynetworks = .*|mynetworks = ${mynetworks}|" /etc/postfix/main.cf
    bashio::log.info "Configure mynetworks: ${mynetworks} - Updated"
else
    echo "mynetworks = ${mynetworks}" >> /etc/postfix/main.cf
    bashio::log.info "Configure mynetworks: ${mynetworks} - Added in config"
fi

bashio::log.info "Init script done!"
